# cron/battery-arhiv.cron.example

*/5 	* * * * 	root 	/etc/voli/cron/battery-arhiv
